# .config/systemd/user/mcps-jupyter-tunnel.service
[Unit]
Description=SSH tunnel to MCPS Jupyter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/ssh -N -L 8889:localhost:8888 mcps
Restart=always
RestartSec=5

# Evita que ssh pida password/passphrase interactivamente y deje el servicio colgado
Environment=AUTOSSH_GATETIME=0
ExecStartPre=/usr/bin/ssh -o BatchMode=yes -o ConnectTimeout=10 mcps true

[Install]
WantedBy=default.target
